# id.tjo.space/root/etc/systemd/system/postgresql-backup.timer
[Unit]
Description=Backup PostgreSQL databases daily

[Timer]
OnCalendar=daily
AccuracySec=1h
Persistent=true
RandomizedDelaySec=100min

[Install]
WantedBy=timers.target
